# net/wifidog/files/wifidog.init
#!/bin/sh /etc/rc.common

START=65
EXTRA_COMMANDS="status"
EXTRA_HELP="        status Print the status of the service"

start() {
    /usr/bin/wifidog-init start
}

stop() {
    /usr/bin/wdctl stop
}

status() {
    /usr/bin/wdctl status
}

restart() {
    stop
    # iptables清理需要时间，所以sleep
    sleep 2
    start
}

reload() {
    restart
}
